# todos.conf
[program:todos]
command=python3 server.py
directory=/pyproject/todos
autostart=true
autorestart=true
